# website_discovery/systemd/url-discovery.service
[Unit]
Description=URL Discovery Service for Bangladeshi Government Websites
Documentation=https://github.com/your-org/website-discovery
After=network.target docker.service
Wants=docker.service

[Service]
Type=simple
User=vpa
Group=vpa
WorkingDirectory=/home/vpa/Prod/url_discovery_engine/website_discovery
EnvironmentFile=/home/vpa/Prod/url_discovery_engine/website_discovery/.env

# Run the main service via conda-activated crawler environment
Environment="PYTHONUNBUFFERED=1"
Environment="PYTHONPATH=/home/vpa/Prod/url_discovery_engine/website_discovery"
ExecStart=/home/vpa/miniconda3/envs/crawler/bin/python -m src.main

# Restart policy - auto-restart on failure
Restart=always
RestartSec=10

# Resource limits
LimitNOFILE=65536
LimitNPROC=256

# Security hardening
NoNewPrivileges=true
ProtectSystem=full
ProtectHome=read-only
PrivateTmp=true
ProtectKernelTunables=true
ProtectKernelModules=true
ProtectControlGroups=true
RestrictRealtime=true
RestrictSUIDSGID=true

# Allow writing to specific directories
ReadWritePaths=/home/vpa/Prod/url_discovery_engine/website_discovery/logs

# Logging
StandardOutput=journal
StandardError=journal
SyslogIdentifier=url-discovery

# Timeout settings
TimeoutStartSec=60
TimeoutStopSec=30

[Install]
WantedBy=multi-user.target
